# supervisord.prod.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start
retries=5
